# debian/postinst
#!/bin/bash
set -e

WORDNET_SYSTEM="/usr/share/wordnet"
# Because postinst is run as root, $HOME evaluates to /root
# Instead, we will place it into /var/lib/wd/wordnet globally.
WORDNET_LOCAL="/var/lib/wd/wordnet"
WORDNET_URL="https://github.com/globalwordnet/english-wordnet/releases/download/2024-edition/english-wordnet-2024.xml.gz"

# Check if system WordNet is installed
if [ -d "$WORDNET_SYSTEM/dict" ]; then
    echo "wd: System WordNet found at $WORDNET_SYSTEM"
else
    echo "wd: System WordNet not found, attempting to download..."
    mkdir -p "$WORDNET_LOCAL"

    # Turn off set -e temporarily so network failures don't break the entire apt install
    set +e
    if command -v wget > /dev/null 2>&1; then
        wget -q -O "/tmp/wordnet-data.gz" "$WORDNET_URL"
        DL_STATUS=$?
    elif command -v curl > /dev/null 2>&1; then
        curl -sL -o "/tmp/wordnet-data.gz" "$WORDNET_URL"
        DL_STATUS=$?
    else
        echo "wd: Neither wget nor curl found."
        DL_STATUS=1
    fi
    set -e

    if [ $DL_STATUS -eq 0 ] && [ -s "/tmp/wordnet-data.gz" ]; then
        echo "wd: WordNet data downloaded. Extracting..."
        # Warning: Your URL points to an .xml.gz, not a tar archive.
        # We use gunzip instead of tar -xzf here.
        gunzip -c "/tmp/wordnet-data.gz" > "$WORDNET_LOCAL/english-wordnet-2024.xml" || true
        rm -f "/tmp/wordnet-data.gz"
        echo "wd: WordNet data installed to $WORDNET_LOCAL"
        
        # Symlink to where the app expects it if needed
        ln -sfn "$WORDNET_LOCAL" /usr/share/wd-wordnet || true
    else
        echo "wd: WARNING: Download failed (Network Error)."
        echo "     Please install wordnet manually: sudo apt install wordnet-base"
        rm -f "/tmp/wordnet-data.gz"
    fi
fi

# Enable systemd user service if user session is active
if [ -n "$DBUS_SESSION_BUS_ADDRESS" ] && command -v systemctl > /dev/null 2>&1; then
    systemctl --user enable --now wd.service 2>/dev/null || true
    echo "wd: Daemon enabled via systemd user service"
else
    echo "wd: Note - run 'systemctl --user enable --now wd.service' after login to start the daemon"
fi

exit 0
